# deploy/discord-video-bot.service
[Unit]
Description=Discord social video rehost bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=discord-bot
Group=discord-bot
WorkingDirectory=/opt/discord-video-bot

# Load secrets and config from this file (not world-readable)
EnvironmentFile=/opt/discord-video-bot/.env

# Prefer a virtualenv if you install one at this path
ExecStart=/opt/discord-video-bot/.venv/bin/python /opt/discord-video-bot/main.py

# Auto-restart on crash or unclean exit
Restart=always
RestartSec=5
StartLimitIntervalSec=0

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/discord-video-bot/logs /tmp/discord-video-bot

# Give the process a clean environment for Python output
StandardOutput=journal
StandardError=journal
SyslogIdentifier=discord-video-bot

[Install]
WantedBy=multi-user.target
